# supervisord-php5.conf
[program:php5]
command=/start-php5.sh
numprocs=1
autostart=true
autorestart=true