# supervisord-httpd.conf
[program:httpd]
command=/start-httpd.sh
numprocs=1
autostart=true
autorestart=true